# .config/systemd/user/codex-poll.service
[Unit]
Description=Codex work-loop poller

[Service]
Type=oneshot
ExecStart=%h/bin/codex-poll
KillMode=process
NoNewPrivileges=true
ProtectHome=read-only
PrivateTmp=true
TimeoutStopSec=30
