# .config/systemd/user/codex-digest.service
[Unit]
Description=Codex founder digest

[Service]
Type=oneshot
ExecStart=%h/scripts/codex_digest.py --once
KillMode=process
